# eink-scheduler.service
# /etc/systemd/system/eink-scheduler.service
#
# Systemd service for E-ink Display Scheduler
# This provides automatic startup and restart on crashes

[Unit]
Description=E-ink Display Scheduler
After=network.target

[Service]
Type=simple
User=pilist
Group=pilist
WorkingDirectory=/home/pilist/eink_scheduler

# Path to your Python script
# I use mamba env
ExecStart=/home/pilist/miniforge3/envs/eink_display_env/bin/python /home/pilist/eink_scheduler/scheduler.py

# Restart policy - automatically restart if it crashes
Restart=always
RestartSec=60

# Standard output and error logging
StandardOutput=journal
StandardError=journal

# Environment variables (if needed)
# Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
